# daemon/com.memorable.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.memorable.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mattkennelly/.pyenv/shims/python3</string>
        <string>/Users/mattkennelly/memorable-plugin/daemon/memorable_daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mattkennelly/memorable-plugin/daemon</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mattkennelly/.memorable/logs/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mattkennelly/.memorable/logs/daemon.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mattkennelly/.local/bin:/Users/mattkennelly/.pyenv/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
